# homelab/ansible/roles/adguard_watchdog/templates/adguard-dns-watchdog.service.j2
# Managed by Ansible — do not edit by hand.
[Unit]
Description=Restart AdGuard if it stops answering DNS
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/adguard-dns-watchdog.sh
