# homelab/ansible/roles/adguard_watchdog/templates/adguard-dns-watchdog.timer.j2
# Managed by Ansible — do not edit by hand.
[Unit]
Description=Run AdGuard DNS watchdog every 30s

[Timer]
OnBootSec=60
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
